# systemd/restic.timer
[Unit]
Description=Daily timer for restic

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
